# Chaskis/Install/linux/systemd/chaskis.service
[Unit]
Description=chaskis
After=network.target

[Service]
Type=forking
User=chaskis
WorkingDirectory=/usr/lib/Chaskis/bin/
ExecStart=/usr/bin/mono-service -l:/tmp/chaskis_$USER.lock -m:chaskis_$USER /usr/lib/Chaskis/bin/ChaskisService.exe
ExecStop=kill `cat /tmp/chaskis_$USER.lock`
TimeoutStopSecs=5
Restart=always

[Install]
Alias=chaskis.service
WantedBy=multi-user.target
