# Chaskis/Install/linux/systemd/chaskis.service
[Unit]
Description=chaskis
After=network.target

[Service]
Type=forking
WorkingDirectory=/usr/lib/Chaskis/bin/
ExecStart=/usr/bin/mono-service -l:/tmp/%u.chaskis.lock -m:%u-chaskis /usr/lib/Chaskis/bin/ChaskisService.exe
ExecStop=/usr/bin/kill `cat /tmp/%u.chaskis.lock`
TimeoutStopSecs=5
Restart=always

[Install]
Alias=chaskis.service
WantedBy=default.target
